# claude_focus.service
[Unit]
Description=Claude Focus Service
After=graphical-session.target

[Service]
Type=simple
ExecStart=/home/tim/.claude/claude_focus_service.py
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-focus

[Install]
WantedBy=default.target