# src/claude_notify/install/templates/claude-notify-daemon.socket
[Unit]
Description=Claude Notify Daemon Socket

[Socket]
ListenStream=%t/claude-notify.sock
SocketMode=0600

[Install]
WantedBy=sockets.target
